# scripts/launchd/com.bjj-betsports.daily-training.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Unique identifier for the job -->
    <key>Label</key>
    <string>com.bjj-betsports.daily-training</string>

    <!-- Script to execute -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__PROJECT_DIR__/scripts/daily_training.sh</string>
    </array>

    <!-- Schedule: every day at 06:00 local time -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/bjj-betsports/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/bjj-betsports/launchd-stderr.log</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <!-- Do not restart on failure (one-shot daily job) -->
    <key>KeepAlive</key>
    <false/>

</dict>
</plist>
